# osu-intercept.service.in
[Unit]
Description=osu-intercept SOCD-cleaning input daemon
Documentation=https://github.com/zakack/osu-intercept
Wants=pipewire.service
After=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/osu-interceptd
Restart=on-failure
RestartSec=2s
Environment=PIPEWIRE_LATENCY=16/48000

# This is a user unit: enable it with
#   systemctl --user enable --now osu-intercept
#
# The daemon needs read access to /dev/input/event* and read/write access
# to /dev/uinput. Both are covered by membership in the `input` group
# (the packaged udev rule opens /dev/uinput to that group):
#   sudo usermod -aG input $USER   # then log out and back in
#
# Config is read from ~/.config/osu-intercept/config.yaml, falling back to
# @CMAKE_INSTALL_FULL_DATADIR@/osu-intercept/config.yaml.

[Install]
WantedBy=default.target
